# navit/navit.service.in
[Unit]
Description=navit software

[Service]
ExecStart="@bindir@/navit"
Restart=on-failure

[Install]
WantedBy=multi-user.target
